# systemd/.config/systemd/user/mpdstatus.service
[Unit]
Description=MPD status daemon

[Service]
Type=simple

ExecStart=/home/luclu7/bins/mpdstatus-daemon

[Install]
WantedBy=default.target
